# nanobsd/plugins/jdownloader_pbi/resources/usr/ports/net/jdownloader/files/jdownloaderd.in
#!/bin/sh
#
# PROVIDE: jdownloader
# REQUIRE: LOGIN 
# KEYWORD: shutdown
# $FreeBSD$
#
# Add the following lines to /etc/rc.conf to enable this service:
#
# jdownloader_enable:            Set to NO by default. Set it to YES to enable it.
#                            a path.
# jdownloader_user:              The user account jdownloader daemon runs as
#

. /etc/rc.subr

name="jdownloader"
rcvar=${name}_enable

pidfile=/var/run/JDownloader/JDownloader.pid
command=%%PREFIX%%/sbin/jdownloader

load_rc_config $name

jdownloader_user=${jdownloader_user="www"}

command_args="$jdownloader_args"

jdownloader_stop()
{
        res=1
        echo -n 'Stopping jdownloader.'
        if [ -f "${pidfile}" ]
        then
                kill -KILL $(cat "${pidfile}")
                rm -f ${pidfile}
                res=$?
        else
                killall -KILL jdownloader
                res=$?
        fi

        echo
        return ${res}
}

stop_cmd='jdownloader_stop'
run_rc_command "$1"

# Need to wait for Java to start before getting PID
sleep 2
pgrep -U www -f JDownloader.jar > /var/run/JDownloader/JDownloader.pid
